# supervisor/conf.d/dnsmasq.ini
[program:dnsmasq]
command=/etc/supervisor/daemon.sh "systemctl start dnsmasq" "/var/run/dnsmasq/dnsmasq.pid"
startsecs=3
redirect_stderr=true
stdout_logfile=/tmp/dnsmasq.log
stdout_logfile_maxbytes=2MB
stdout_logfile_backups=1
