# supervisor/conf.d/if.ini
[program:wlan-daemon]
command=/root/src/arm/interface/if_keep.sh wlan1 wlan3
startsecs=10
redirect_stderr=true
stdout_logfile=/tmp/wlan-daemon.log
stdout_logfile_maxbytes=2MB
stdout_logfile_backups=1
