# supervisor/conf.d/nat.ini
[program:nat-reset]
directory=/root/src/arm/iptables
command=/root/src/arm/iptables/reset.sh
autostart=false
startsecs=1
startretries=0
redirect_stderr=true
stdout_logfile=/tmp/nat-reset.log
stdout_logfile_maxbytes=2MB
stdout_logfile_backups=1

[program:nat-ap]
directory=/root/src/arm/iptables
command=/root/src/arm/iptables/rules/ap.sh
autostart=false
startsecs=1
startretries=0
redirect_stderr=true
stdout_logfile=/tmp/nat-ap.log
stdout_logfile_maxbytes=2MB
stdout_logfile_backups=1

[program:nat-ss]
directory=/root/src/arm/iptables
command=/root/src/arm/iptables/rules/ss-nat.sh install
autostart=false
startsecs=1
startretries=0
redirect_stderr=true
stdout_logfile=/tmp/nat-ss.log
stdout_logfile_maxbytes=2MB
stdout_logfile_backups=1

[program:nat-ss-remove]
directory=/root/src/arm/iptables
command=/root/src/arm/iptables/rules/ss-nat.sh remove
autostart=false
startsecs=1
startretries=0
redirect_stderr=true
stdout_logfile=/tmp/nat-ss.log
stdout_logfile_maxbytes=2MB
stdout_logfile_backups=1
