# supervisor/conf.d/pen.ini
[program:dns-proxy]
command=pen -dfrU -S 3 20000 127.0.0.1:1053 127.0.0.1:1054 127.0.0.1:1055
startsecs=2
startretries=100
redirect_stderr=true
stdout_logfile=/tmp/dns-proxy.log
stdout_logfile_maxbytes=2MB
stdout_logfile_backups=1
