# supervisor/conf.d/ss.ini
[program:ss-redir]
command=ss-redir -c /etc/ss/ss.json  -l 1080 -vvv
startsecs=10
redirect_stderr=true
stdout_logfile=/tmp/ss-redir.log
stdout_logfile_maxbytes=2MB
stdout_logfile_backups=1

[program:ss-tunnel-01]
command=ss-tunnel -c /etc/ss/ss-tunnel.json -l 1053 -L 8.8.8.8:53 -u -vvv
startsecs=10
redirect_stderr=true
stdout_logfile=/tmp/ss-tunnel-01.log
stdout_logfile_maxbytes=2MB
stdout_logfile_backups=1

[program:ss-tunnel-02]
command=ss-tunnel -c /etc/ss/ss-tunnel.json -l 1054 -L 8.8.8.8:53 -u -vvv
startsecs=10
redirect_stderr=true
stdout_logfile=/tmp/ss-tunnel-02.log
stdout_logfile_maxbytes=2MB
stdout_logfile_backups=1

[program:ss-tunnel-03]
command=ss-tunnel -c /etc/ss/ss-tunnel.json -l 1055 -L 8.8.8.8:53 -u -vvv
startsecs=10
redirect_stderr=true
stdout_logfile=/tmp/ss-tunnel-03.log
stdout_logfile_maxbytes=2MB
stdout_logfile_backups=1

[group:ss]
programs=ss-redir,ss-tunnel-01,ss-tunnel-02,ss-tunnel-03 
priority=999