# supervisor/conf.d/udhcpd.ini
[program:udhcpd]
command=/etc/supervisor/daemon.sh "/usr/sbin/udhcpd -S" "/var/run/udhcpd.pid"
startsecs=3
redirect_stderr=true
stdout_logfile=/tmp/udhcpd.log
stdout_logfile_maxbytes=2MB
stdout_logfile_backups=1